# deploy/cuchuma.service
[Unit]
Description=Portal de Defensa y Documentación de El Cuchumá
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/cuchuma
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/cuchuma/.env
ExecStart=/usr/bin/node build
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
